# internal/templates/launchd/com.gastown.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gastown.daemon</string>

    <!-- Sources <town-root>/settings/daemon.env (if present, mode 0600) so
         secrets like GT_TELEGRAPH_*_SECRET reach the daemon. launchd inherits
         a stripped environment, so values from ~/.zshrc don't propagate.
         `set -a` makes plain `KEY=value` lines auto-export, matching systemd's
         EnvironmentFile semantics so the same file works on both supervisors. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>set -a; [ -r "{{.TownRoot}}/settings/daemon.env" ] &amp;&amp; . "{{.TownRoot}}/settings/daemon.env"; set +a; exec "{{.GTPath}}" daemon run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{.TownRoot}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>{{.TownRoot}}/daemon/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>{{.TownRoot}}/daemon/daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GT_TOWN_ROOT</key>
        <string>{{.TownRoot}}</string>
        <key>PATH</key>
        <string>{{.Path}}</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
